# deploy/droplet/phase0/systemd/hermes-notify-queue.timer
[Unit]
Description=Deliver queued Hermes alerts every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=hermes-notify-queue.service

[Install]
WantedBy=timers.target
